# files/mqtt/mqtt.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/mqtt.sh
User=lbischof
Group=lbischof

[Install]
WantedBy=default.target
