# files/udev/rules.d/99-monitor-hotplug.rules
KERNEL=="card0", SUBSYSTEM=="drm", RUN+="/home/lbischof/monitor-hotplug.sh"
