# roles/desktop/files/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/agetty --autologin lbischof --noclear %I $TERM
Type=simple
